# systemd/buscount-ramdisk.mount.in
[Unit]
Description=Buscount ramdisk

[Mount]
What=tmpfs
Where=@RAMDISK_LOCATION@
Type=tmpfs
Options=size=50M

[Install]
WantedBy=multi-user.target

