# systemd/buscountd.service.in
[Unit]
Description=Buscountd service
After=network.target

[Service]
ExecStartPre=+@NCS_RESET@
ExecStart=bash -c ". @OPENVINO_SETUPVARS@; @BUSCOUNTD_LOCATION@"
User=pi
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
Wants=@RAMDISK_NAME@

